# laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=/opt/plesk/php/7.3/bin/php /var/www/vhosts/smart-edu.id/api.smart-edu.my.id/artisan queue:work sync --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/vhosts/smart-edu.id/api.smart-edu.my.id/worker.log
stopwaitsecs=3600
